# data/data-page/data-detection/list-detection.json
[
      { "Detection-ID": "DET-1", "Detection Name": "Inhibit System Recovery", "Mitre Attack ID": "T1490", "Comment": "N/A", "MITRE CAR": "CAR-2021-01-009"  },
      { "Detection-ID": "DET-2", "Detection Name": "Indicator Removal on Host: File Deletion", "Mitre Attack ID": "T1070.004", "Comment": "N/A" },
      { "Detection-ID": "DET-3", "Detection Name": "Software Discovery: Security Software Discovery", "Mitre Attack ID": "T1518.001", "Comment": "N/A" },
      { "Detection-ID": "DET-4", "Detection Name": "Impair Defenses: Disable or Modify Tools", "Mitre Attack ID": "T1562.001", "Comment": "N/A" },
      { "Detection-ID": "DET-5", "Detection Name": "Domain Trust Discovery", "Mitre Attack ID": "T1482", "Comment": "N/A" },
      { "Detection-ID": "DET-6", "Detection Name": "Permission Groups Discovery: Domain Groups", "Mitre Attack ID": "T1069.002", "Comment": "N/A" },
      { "Detection-ID": "DET-7", "Detection Name": "System Information Discovery", "Mitre Attack ID": "T1426", "Comment": "N/A" },

      { "Detection-ID": "DET-8", "Detection Name": "System Time Discovery", "Mitre Attack ID": "T1124", "Comment": "N/A" },
      { "Detection-ID": "DET-9", "Detection Name": "Hide Artifacts: Hidden Files and Directories", "Mitre Attack ID": "T1564.001", "Comment": "N/A" },
      { "Detection-ID": "DET-10", "Detection Name": "Password Policy Discovery", "Mitre Attack ID": "T1201", "Comment": "N/A" },
      { "Detection-ID": "DET-11", "Detection Name": "Network Share Creation", "Mitre Attack ID": "N/A", "Comment": "Share Creation" },
      { "Detection-ID": "DET-12", "Detection Name": "Account Manipulation", "Mitre Attack ID": "T1098", "Comment": "Password modification" },
      { "Detection-ID": "DET-13", "Detection Name": "File and Directory Permissions Modification: Windows File and Directory Permissions Modification", "Mitre Attack ID": "T1222.001", "Comment": "N/A" },
      { "Detection-ID": "DET-14", "Detection Name": "Masquerading: Match Legitimate Name or Location", "Mitre Attack ID": "T1036.005", "Comment": "N/A" },
      { "Detection-ID": "DET-15", "Detection Name": "Command and Scripting Interpreter", "Mitre Attack ID": "T1059", "Comment": "UsageLogs File Creation", "Forensic": "Yes"  },
      { "Detection-ID": "DET-16", "Detection Name": "Malicious Executable Creation", "Mitre Attack ID": "N/A", "Comment": "Malicious Executable Creation" },
      { "Detection-ID": "DET-17", "Detection Name": "Scheduled Task/Job: Scheduled Task", "Mitre Attack ID": "T1053.005", "Comment": "N/A" },
      { "Detection-ID": "DET-18", "Detection Name": "Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder", "Mitre Attack ID": "T1547.001", "Comment": "N/A" },
      { "Detection-ID": "DET-19", "Detection Name": "Network query", "Mitre Attack ID": "N/A", "Comment": "Ping Query, DNS Query" },
      { "Detection-ID": "DET-20", "Detection Name": "Modify Registry", "Mitre Attack ID": "T1112", "Comment": "N/A" },
      { "Detection-ID": "DET-21", "Detection Name": "Command and Scripting Interpreter: Visual Basic", "Mitre Attack ID": "T1059.005", "Comment": "N/A" },
      { "Detection-ID": "DET-22", "Detection Name": "Deobfuscate Decode Files or Information", "Mitre Attack ID": "T1140", "Comment": "N/A" },


      { "Detection-ID": "DET-XXX", "Detection Name": "Possible malicious pattern", "Mitre Attack ID": "N/A", "Comment": "Possible malicious pattern" }
]
